# Scripts/Install/CWF/etc/xinetd.d/rghttp-service
service vista-rghttp
{
  port = 9080
  socket_type = stream
  protocol = tcp
  type = UNLISTED
  user = osehra
  server = /home/osehra/bin/rghttp.sh
  wait = no
  disable = no
  per_source = UNLIMITED
  instances = UNLIMITED
}
